# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:sidekiq]
command=bundle exec sidekiq -C config/sidekiq.yml
directory=/rails
